# profile.d/java.sh
PATH=$PATH:/usr/java/lasest/bin
export PATH
